# src/lib/systemd/network-sniffer.service
[Unit]
Description=Network Sniffer Suite
After=network-online.target

[Service]
Type=simple
User=network-sniffer
Group=network-sniffer
WorkingDirectory=/tmp

# Seguridad de archivos
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# Importante: No pongas la misma ruta en ReadWrite y ReadOnly
ReadWritePaths=/tmp/ /var/log/network-sniffer
ReadOnlyPaths=/etc/network-sniffer /usr/bin/ss /usr/sbin/ss

# Capacidades
# ss -p necesita "inspeccionar" procesos. En kernels modernos, 
# CAP_SYS_PTRACE es a veces necesario junto a DAC_READ_SEARCH para ver FDs de otros.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN CAP_DAC_READ_SEARCH CAP_SYS_PTRACE
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN CAP_DAC_READ_SEARCH CAP_SYS_PTRACE

# Visibilidad de procesos
# Cambiamos ProtectProc a 'invisible' o 'default' pero asegurando acceso a /proc
ProtectProc=default
ProcSubset=all

# Protocolos
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_PACKET AF_UNIX
NoNewPrivileges=yes

ExecStart=/bin/bash /etc/network-sniffer/network-sniffer.sh

Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target